# systemd/raspi-sentinel-tmpfs-verify.service
[Unit]
Description=Verify raspi-sentinel tmpfs storage is writable and safe
After=run-raspi\x2dsentinel.mount
Before=raspi-sentinel.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/raspi-sentinel/.venv/bin/raspi-sentinel -c /etc/raspi-sentinel/config.toml verify-storage --expected-mode 0755 --expected-owner-uid 0 --expected-owner-gid 0
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
